# AgentKVTMac/Deploy/com.agentkvt.macapp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">

<!--
  AgentKVT Mac App — launchd User Agent
  =====================================
  Use this when the production/TestFlight macOS app is installed and you want the
  signed app bundle to relaunch automatically after login or unexpected exit.

  Install:
    cp com.agentkvt.macapp.plist ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.agentkvt.macapp.plist
    launchctl enable gui/$(id -u)/com.agentkvt.macapp

  Uninstall:
    launchctl bootout gui/$(id -u) ~/Library/LaunchAgents/com.agentkvt.macapp.plist
    rm ~/Library/LaunchAgents/com.agentkvt.macapp.plist

  Notes:
  - Replace the app path below with the actual installed app location.
  - We launch via `open -W -a` (omit `-n`) so a relaunch does not start a second
    instance while one is already running; LaunchServices starts a real macOS app
    and launchd waits until that app exits.
  - The app reads its runtime settings from:
      ~/Library/Group Containers/group.com.agentkvt.shared/Library/Application Support/agentkvt-runner.plist
    so you usually do not need to pass environment variables here.
  - For an unsigned or Xcode-built app launched over SSH, app-group file access
    may hang before the scheduler starts. In that case, add LaunchAgent
    environment variables for:
      AGENTKVT_CONFIG_FILE=~/.agentkvt/agentkvt-runner.plist
      AGENTKVT_LOG_FILE=~/.agentkvt/logs/agentkvt-macapp.log
    and copy the runtime plist to the home-path config file.
  - Runtime logs live at:
      ~/Library/Group Containers/group.com.agentkvt.shared/Library/Logs/agentkvt-mac.log
  - If you need to change OLLAMA_MODEL in the app-group plist over SSH: do not run
    PlistBuddy as root against another user's group container (often permission /
    TCC issues). Use repo script set-ollama-model-in-group-plist.sh as admin, or run
    PlistBuddy/python as that user, e.g. sudo -u familyagent python3 ... on the plist path.
  - launchctl bootout/bootstrap must use the same account that owns the LaunchAgent
    (e.g. sudo -u familyagent launchctl ...), not your personal GUI uid.

  Security — do NOT put API keys or tokens in this LaunchAgent plist. Files under
  ~/Library/LaunchAgents are often copied into diagnostics or dotfile repos. Store
  OLLAMA_API_KEY, AGENTKVT_AGENT_TOKEN, etc. in ~/.agentkvt/agentkvt-runner.plist
  (chmod 600) and point the app there with AGENTKVT_CONFIG_FILE if needed. Use:
    bin/agentkvt_move_ollama_key_off_launchagent.sh
-->

<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agentkvt.macapp</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/open</string>
        <string>-W</string>
        <string>-a</string>
        <string>/Applications/AgentKVTMacApp.app</string>
    </array>

    <key>LimitLoadToSessionType</key>
    <array>
        <string>Aqua</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>60</integer>

    <key>StandardOutPath</key>
    <string>/tmp/agentkvt-macapp-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/agentkvt-macapp-launchd.err</string>
</dict>
</plist>
